# deploy/bikebeach.service

# Systemd service example for BikeBeach (Gunicorn)
# This unit uses an EnvironmentFile so you can set deployment-specific values
# in `/etc/default/bikebeach` (or edit the Environment= lines below).
# Example `/etc/default/bikebeach`:
# WORKDIR=/srv/bikebeach
# VENV_BIN=/srv/bikebeach/.venv/bin
# USER=www-data
# GROUP=www-data

[Unit]
Description=Gunicorn instance to serve BikeBeach
After=network.target

[Service]
# Load variables from file if present
EnvironmentFile=-/etc/default/bikebeach
# Fallback defaults (override in /etc/default/bikebeach)
Environment="WORKDIR=/srv/bikebeach"
Environment="VENV_BIN=/srv/bikebeach/.venv/bin"
Environment="USER=www-data"
Environment="GROUP=www-data"

User=%p
Group=%p
WorkingDirectory=${WORKDIR}

# Use the virtualenv gunicorn binary
ExecStart=${VENV_BIN}/gunicorn --access-logfile - --workers 3 --bind 127.0.0.1:8000 bikebeach.wsgi:application

[Install]
WantedBy=multi-user.target
